# tarpyt@.socket.in
[Unit]
Description=An @PROTOCOL@ tarpit on port %i
After=network.target

[Socket]
ListenStream=%i

[Install]
WantedBy=sockets.target