# tarpyt@.socket.in
[Unit]
Description=An @PROTOCOL@ tarpit on port %i
After=network.target

[Socket]
ListenStream=%i

# Turn on IP accounting
IPAccounting=yes

[Install]
WantedBy=sockets.target
